# roles/consul/templates/consul.socket
[Unit]
Description=dns socket

[Socket]
ListenStream=53
NoDelay=true